# mnheme.service
# /etc/systemd/system/mnheme.service
[Service]
WorkingDirectory=/opt/mnheme
ExecStart=uvicorn mnheme_api:app --host 127.0.0.1 --port 8000
Restart=always
